# server1-supervisord.conf
[program:daemon_print]
command=/bin/bash /tmp/daemon_print
stdout_logfile=/tmp/daemon_print.log
stderr_logfile=/tmp/daemon_print.err
user=training
